# deploy/supervisor.conf
; 宝塔 Supervisor 配置模板
; 使用前将 /www/wwwroot/zhengtongrenhe 替换为你的实际路径
; 在宝塔「Supervisor 管理器」中分别添加，或合并导入

[program:edu-policy-api]
command=/www/wwwroot/zhengtongrenhe/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
directory=/www/wwwroot/zhengtongrenhe/backend
user=www
autostart=true
autorestart=true
stdout_logfile=/www/wwwlogs/edu-policy-api.log
stderr_logfile=/www/wwwlogs/edu-policy-api.err.log

[program:edu-policy-celery]
command=/www/wwwroot/zhengtongrenhe/backend/.venv/bin/celery -A app.tasks.celery_app worker --loglevel=info --concurrency=1
directory=/www/wwwroot/zhengtongrenhe/backend
user=www
autostart=true
autorestart=true
stdout_logfile=/www/wwwlogs/edu-policy-celery.log
stderr_logfile=/www/wwwlogs/edu-policy-celery.err.log

[program:edu-policy-beat]
command=/www/wwwroot/zhengtongrenhe/backend/.venv/bin/celery -A app.tasks.celery_app beat --loglevel=info
directory=/www/wwwroot/zhengtongrenhe/backend
user=www
autostart=true
autorestart=true
stdout_logfile=/www/wwwlogs/edu-policy-beat.log
stderr_logfile=/www/wwwlogs/edu-policy-beat.err.log
